# multitest_transport/scripts/crontab
# A crontab file for Docker container.
0 * * * * /mtt/scripts/file_cleaner.sh >> /data/log/file_cleaner.log
